# crontab
*/1 * * * * python3 /selen/Selenium.py >> /var/log/cron.log 2>&1
